# ansible/roles/proxmox_host/templates/management-nic-watchdog.service.j2
[Unit]
Description=Check {{ management_interface }} for e1000e hardware hangs
Wants=network-online.target
After=network-online.target
ConditionPathExists=/sys/class/net/{{ management_interface }}

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/management-nic-watchdog
